# config/systemd/alarm@.service
[Unit]
Description=One-shot alarm (%i)

[Service]
ExecStart=%h/scripts/alarm/trigger.sh %i
