# config/systemd/timer@.service
[Unit]
Description=Timer service (%i seconds)

[Service]
Type=oneshot
ExecStart=%h/scripts/timer/trigger.sh %i
