# config/systemd/timer@.timer
[Unit]
Description=Timer (%i seconds)

[Timer]
OnActiveSec=%i
AccuracySec=1s
Unit=timer@%i.service

[Install]
WantedBy=timers.target
